# deployment/djangotcha.conf
description "Gunicorn DjanGotcha"

start on runlevel [2345]
stop on runlevel [06]

exec su -c "/home/ubuntu/git/djangotcha/deployment/djangotcha.sh"
respawn